# deploy/systemd/agent-gateway.service
[Unit]
Description=Agent Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/agent-gateway.env
WorkingDirectory=/
ExecStart=/bin/bash -lc 'cd "$BRIDGE_ROOT" && exec ./scripts/restart-supervisor.sh -- "$NODE_BIN" ./scripts/start-with-proxy.mjs'
Restart=always
RestartSec=3
TimeoutStopSec=20
KillMode=mixed

[Install]
WantedBy=multi-user.target
